# systemd/seedcare-dashboard.service
[Unit]
Description=SeedCare Dashboard (Streamlit)
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/app/seedcare
EnvironmentFile=-/opt/app/seedcare/.env
ExecStart=/opt/app/seedcare/.venv/bin/streamlit run seedcare/dashboard.py \
    --server.port 8501 --server.address 0.0.0.0 --server.headless true \
    --server.enableCORS false --server.enableXsrfProtection false \
    --server.enableWebsocketCompression false \
    --browser.serverAddress 0.0.0.0
Restart=on-failure
RestartSec=5
Environment=PATH=/opt/app/seedcare/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/app/seedcare

[Install]
WantedBy=multi-user.target
